# debian/DEBIAN/postinst
#!/bin/bash

echo "\033[36m adsbxfeed postinst `date +"%s"` \033[39m"

systemctl enable adsbxfeed
systemctl start adsbxfeed

exit 0