# debian/DEBIAN/prerm
#!/bin/bash

echo "\033[36m adsbxfeed prerm `date +"%s"` \033[39m"

systemctl disable adsbxfeed
systemctl stop adsbxfeed

exit 0