# debian/DEBIAN/prerm
#!/bin/bash

echo "adsbxfeed prerm"

systemctl disable adsbxfeed
systemctl stop adsbxfeed

deluser --system adsbxfeed

exit 0